# Library/LaunchAgents/com.henry.claude-window-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.henry.claude-window-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>WATCHER_SCRIPT_PATH</string>
    </array>

    <!-- Keep alive always — the script has its own sleep/loop logic -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start immediately on load -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/tsaishengwei/.claude/window-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tsaishengwei/.claude/window-watcher.log</string>

    <!-- Throttle restarts: wait 60s before relaunching if script crashes -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
